# secor/kafka/supervisor/mqtt.conf
[program:mqtt]
command=/usr/sbin/mosquitto
autostart=true
autorestart=true
