# secor/kafka/supervisor/mqtt2kafka.conf
[program:mqtt2kafka]
command=/usr/bin/MQTTtoKAFKA.sh
autostart=true
autorestart=true
