# secor/kafka/supervisor/rest-proxy.conf
[program:kafka-rest-proxy]
command=/usr/bin/rest-proxy-docker.sh
autostart=true
autorestart=true
