# secor/SECOR_INSTALL_DIR/supervisor/secor.conf
[program:secor]
command=/opt/secor/secor_start.sh
autostart=true
autorestart=true
